# systemd/home-mocha-Documents-Dropbox.automount
[Unit]
Description=Rclone automount for Dropbox documents

[Automount]
Where=/home/mocha/Documents/Dropbox/
TimeoutIdleSec=600

[Install]
WantedBy=default.target
