# systemd/rclone-Pictures.timer
[Unit]
Description="Bisync Pictures every 5 minutes"

[Timer]
OnBootSec=10sec
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

